# multi-instance/templates/openclaw-instance.container.tmpl
# OpenClaw {{NAME}} instance - Podman Quadlet (rootless, hardened)
# Generated by openclaw-instance.sh
# Security: ports mapped to 127.0.0.1 only on host

[Unit]
Description=OpenClaw {{NAME}} gateway (rootless Podman)
After=network-online.target
Wants=network-online.target

[Container]
Image={{IMAGE}}
ContainerName=openclaw-{{NAME}}
UserNS=keep-id

# Loopback-only binding on HOST (security: no LAN exposure)
PublishPort=127.0.0.1:{{GATEWAY_PORT}}:{{GATEWAY_PORT}}
PublishPort=127.0.0.1:{{BRIDGE_PORT}}:{{BRIDGE_PORT}}

# Mount config, workspace, and shared skills
Volume={{OPENCLAW_HOME}}/.openclaw-{{NAME}}:/home/node/.openclaw:rw
Volume={{OPENCLAW_HOME}}/workspace-{{NAME}}:/home/node/workspace:rw
Volume={{OPENCLAW_HOME}}/shared/skills:/home/node/shared-skills:ro

# Environment
EnvironmentFile={{OPENCLAW_HOME}}/.openclaw-{{NAME}}/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color

# DNS (MagicDNS + fallback)
DNS={{DNS_PRIMARY}}
DNS={{DNS_FALLBACK}}

# Resource limits + user mapping (must match host uid:gid)
PodmanArgs=--memory={{MEMORY}} --cpus={{CPUS}} --user {{USER_MAPPING}}

Pull=never
# Use --bind lan inside container; host port mapping enforces loopback
Exec=node dist/index.js gateway --bind lan --port {{GATEWAY_PORT}}

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
